# quick-tunnel/wedding-redirect.timer
[Unit]
Description=Check + publish wedding tunnel URL every minute

[Timer]
OnBootSec=30
OnUnitActiveSec=60
AccuracySec=10

[Install]
WantedBy=timers.target
